# lib/premades/RegisterTask.xml
<?xml version="1.0" ?>
<!-- RegisterTask.xml
 This project, which includes this file, is licensed under GNU General
 Public License v3.0.
 Get a copy here: https://www.gnu.org/licenses/gpl-3.0-standalone.html
 Or just look at the LICENSE file.
 Last Updated 16 June 2021 -->

<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
    <RegistrationInfo>
        <Date>2021-06-11T8:00:00-05:00</Date>
        <Author>ikeacat (NoPerish)</Author>
        <Version>1.0.0</Version>
        <Description>Starts the NoPerish startup script on logon.</Description>
    </RegistrationInfo>
    <Triggers>
        <LogonTrigger>
            <StartBoundary>2021-06-11T8:00:00-05:00</StartBoundary>
            <EndBoundary>2999-12-31T12:59:59</EndBoundary>
            <Enabled>true</Enabled>
            <UserId>UserName</UserId>
        </LogonTrigger>
    </Triggers>
    <Principals>
        <Principal>
            <GroupId>Builtin\Administrators</GroupId>
        </Principal>
    </Principals>
    <Settings>
        <Enabled>true</Enabled>
        <AllowStartOnDemand>false</AllowStartOnDemand>
        <AllowHardTerminate>true</AllowHardTerminate>
    </Settings>
    <Actions>
        <Exec>
            <Command>start %USERPROFILE%/NoPerish/noperish</Command>
        </Exec>
    </Actions>
</Task>